# rootfs/etc/services.d/mayan_celery_beat/run
#!/usr/bin/with-contenv bash

nice -n 1 /bin/bash -c "${MAYAN_BIN} celery --settings=${MAYAN_SETTINGS_MODULE} beat --pidfile= -l ERROR"

